# init/papasan.conf
description "papasan service"

start on stopped rc RUNLEVEL=[2345]

respawn
respawn limit 5 60
oom never

limit stack unlimited unlimited

exec /usr/sbin/papasan
